# s6-overlay/s6-rc.d/generate-api-key-and-secret/up
/etc/s6-overlay/scripts/generate-api-key-and-secret 2>&1