# s6-overlay/s6-rc.d/torrent-archiver/run
#!/command/with-contenv sh
set -a
source /etc/webtor/common.env
WEB_PORT=$TORRENT_ARCHIVER_SERVICE_PORT
set +a
/app/torrent-archiver 2>&1 | s6-log p"[torrent-archiver]" 1